# profile.d/80-theme.el
;;; -*- lexical-binding: t; -*-
;;; 80-theme.el --- Theme profile
;;; Commentary:

;;; Code:

;; (load-theme 'spacemacs-dark)

;;; 80-theme.el ends here
;;; End:
