# profile.d/90-custom.el
;;; 90-custom.el --- Custom profile configuration
;;; Commentary:

;;; Code:

(use-package ldap-mode)

;;; 90-custom.el ends here
;;; End:
