# profile.d/99-server.el
;;; 99-server.el --- server configuration
;;; Commentary:

;;; Code:

;;(server-start)

;;; 99-server.el ends here
;;; End:
